# etc/services.d/webapp/run
#!/usr/bin/execlineb -P
exec webapp_skeleton

